# etc/s6/ssh/run
#!/bin/sh

exec sh -c "echo 'PermitRootLogin yes' >> /etc/ssh/sshd_config && ssh-keygen -f /etc/ssh/ssh_host_rsa_key -N '' -t rsa && /usr/sbin/sshd -D -f /etc/ssh/sshd_config"
