# etc/s6/consul-template/run
#!/bin/sh

exec /usr/local/bin/consul-template -config /etc/consul-template.d
